# drift-beacon/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::net.wait_for $SERVICE_PORT localhost 900

bashio::log.info "Starting nginx..."
exec nginx
